# profile.d/vim_alias.sh
# For clipboard support
alias vim='gvim -v'
